# usr/share/netrunner-default-settings/plasma5-profile/autostart/novnc.desktop
[Desktop Entry]
Exec=websockify --web /usr/share/novnc/ 5901 localhost:5900
Terminal=false
Type=Application
X-DBUS-StartupType=unique
