# usr/share/netrunner-default-settings/plasma5-profile/autostart/vnc.desktop
[Desktop Entry]
Exec=websockify --wrap-mode=respawn --web /usr/share/novnc/ 5900 -- krfb
Terminal=false
Type=Application
X-DBUS-StartupType=unique
