# usr/share/netrunner-default-settings/plasma5-profile/autostart/vnc.desktop
[Desktop Entry]
Exec=x11vnc -forever -shared -xdmsvc -ncache_cr -http -http_oneport
Terminal=false
Type=Application
X-DBUS-StartupType=unique
