# usr/share/netrunner-default-settings/plasma5-profile/autostart/vnc.desktop
[Desktop Entry]
Exec=x11vnc -norepeat -shared -xdmsvc -loop
Terminal=false
Type=Application
X-DBUS-StartupType=unique
